# resumeforge-backend.service
[Unit]
Description=ResumeForge Backend API
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=johngalt
ExecStart=/var/www/resumeforge.davidveksler.com/ResumeForge/venv/bin/python /var/www/resumeforge.davidveksler.com/ResumeForge/app.py
WorkingDirectory=/var/www/resumeforge.davidveksler.com/ResumeForge
Environment=PATH=/var/www/resumeforge.davidveksler.com/ResumeForge/venv/bin
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=False
Environment=PYTHONPATH=/var/www/resumeforge.davidveksler.com/ResumeForge
Environment=PORT=5000
EnvironmentFile=-/var/www/resumeforge.davidveksler.com/ResumeForge/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=resumeforge-backend

[Install]
WantedBy=multi-user.target